# mcp-services/com.mcp.ssh-manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mcp.ssh-manager</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npx</string>
    <string>-y</string>
    <string>supergateway</string>
    <string>--stdio</string>
    <string>npx -y mcp-ssh-manager</string>
    <string>--port</string>
    <string>8001</string>
    <string>--healthEndpoint</string>
    <string>/healthz</string>
  </array>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/YOUR_USER/.mcp-services/logs/ssh-manager.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USER/.mcp-services/logs/ssh-manager.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
