# mcp-services/com.mcp.outlook-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mcp.outlook-mcp</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npx</string>
    <string>-y</string>
    <string>supergateway</string>
    <string>--stdio</string>
    <string>npx outlook-mcp</string>
    <string>--port</string>
    <string>8003</string>
    <string>--healthEndpoint</string>
    <string>/healthz</string>
  </array>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/YOUR_USER/.mcp-services/logs/outlook-mcp.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USER/.mcp-services/logs/outlook-mcp.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>MS_CLIENT_ID</key>
    <string>YOUR_MS_CLIENT_ID</string>
    <key>MS_CLIENT_SECRET</key>
    <string>YOUR_MS_CLIENT_SECRET</string>
    <key>MS_TENANT_ID</key>
    <string>YOUR_MS_TENANT_ID</string>
  </dict>
</dict>
</plist>
